# config/dpatch.supervisor.ini
[program:tasks]
numprocs = 1
user = dpatch
process_name = dailytasks-%(process_num)s
directory = /home/dpatch
command = /usr/share/dpatch/bin/dailytasks.sh
environment=HOME="/home/dpatch",USER="dpatch"
autorestart = true
startsecs = 60
redirect_stderr = true
stdout_logfile = /var/log/dpatch/dailytasks.log
stderr_logfile = /var/log/dpatch/dailytasks-error.log
